# voyant-export/ai-cli-crash-capture/systemd/terminal-recoverer-coredump.service
[Unit]
Description=Terminal recoverer coredump watcher (best-effort)
After=default.target

[Service]
ExecStart=%h/.local/bin/terminal-recoverer watch-coredump
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
